# autostart/adb-media-player.desktop
[Desktop Entry]
Name=ADB Media Player
Exec=/home/sled/react-car-media-player/dist/adb-media-player/adb-media-player-linux_arm64
Type=Application
