# autostart/android-controller.desktop
[Desktop Entry]
Name=Android Controller
Exec=/home/sled/react-car-media-player/lua/android.lua
Type=Application
