# autostart/remove-lxpanel.desktop
[Desktop Entry]
Name=File Manager
Exec=pkill lxpanel
Type=Application
